# deploy/linux/systemd/equinox-backup.service
[Unit]
Description=Backup diario de Equinox Gym
After=mysql.service

[Service]
Type=oneshot
User=root
ExecStart=/opt/equinox-gym/backup.sh
